# scripts/systemd/research-bot-slack.service
[Service]
Type=simple
WorkingDirectory=%h/research-bot
EnvironmentFile=%h/research-bot/.env
ExecStart=%h/.cargo/bin/cargo run --bin slack_bot
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
