# scripts/systemd/research-bot-rag-sync.service
[Service]
Type=oneshot
EnvironmentFile=%h/research-bot/rag-sync.env
ExecStart=%h/research-bot/scripts/sync_rag_artifacts.sh
