# scripts/systemd/research-bot-slack.service
[Unit]
Description=Research Bot Slack (Socket Mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/research-bot
EnvironmentFile=/etc/research-bot/research-bot.env
ExecStart=/usr/bin/cargo run --bin slack_bot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
