# scripts/systemd/research-bot-rag-sync.service
[Unit]
Description=Sync RAG artifacts for research-bot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/research-bot/rag-sync.env
ExecStart=/opt/research-bot/scripts/sync_rag_artifacts.sh
